# root/etc/services.d/rutorrent/run
#!/usr/bin/with-contenv bash
umask 002

exec \
	screen -d -m -S \
	rtorrent s6-setuidgid abc /usr/bin/rtorrent \
	-n -o import=/config/rtorrent/rtorrent.rc

until [ -e "/rtorrent/rtorrent_sess/rtorrent.lock" ];
do
sleep 1s
done

rtorrent_pid=$(cat /rtorrent/rtorrent_sess/rtorrent.lock | cut -d '+' -f 2)
tail -n 1 -f /config/log/rtorrent/rtorrent.log "$rtorrent_pid"

